# robot/rover/systemd/arm-command-listener.service
[Unit]
Description=Arm command listener

[Service]
ExecStart=/usr/bin/ArmCommandListener.py
